# deploy/skzdev02/premarket-backend.service
[Unit]
Description=Pre-Market Dashboard Backend
After=network.target

[Service]
Type=simple
User=adminskz
WorkingDirectory=/opt/premarket/backend
EnvironmentFile=/opt/premarket/backend/.env
# Bind to localhost only; nginx (port 3000) reverse-proxies /api to it.
ExecStart=/opt/premarket/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal

# Mild hardening. NOTE: do NOT add ReadOnlyPaths=/ — the app writes the SQLite
# DB under /opt/premarket/data, and that would break it.
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
